# 1.10/build/configs/supervisord/ost.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.2 -F
stopasgroup=true
stopsignal=TERM
autostart=true
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stderr_logfile_maxbytes=5MB
stdout_logfile=/var/log/supervisor/php-fpm.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=7


[program:cron]
command=/usr/sbin/cron -f
stopsignal=TERM
stopasgroup=true
autostart=true
stderr_logfile=/var/log/supervisor/cron.err.log
stderr_logfile_maxbytes=5MB
stdout_logfile=/var/log/supervisor/cron.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=7
